# facebook/supervisor/conf.d/campaign.conf
[program:campaign]
directory=/app
command=python3 /app/campaign.py
autostart=true
autorestart=true
redirect_stderr=true
